# greyhole.cron.weekly
#!/bin/sh

# Weekly fsck
/usr/bin/greyhole --fsck --email-report --dont-walk-graveyard > /dev/null

