# greyhole.cron.d
0 0 * * *	root	/usr/bin/greyhole-executer --fsck --email-report > /dev/null