# greyhole.cron.daily
#!/bin/sh

# Daily (conditional) fsck
/usr/bin/greyhole --fsck --if-conf-changed --email-report --dont-walk-metadata-store > /dev/null

