# upstart.conf
description "Greyhole - Easily expandable & redundant storage pool using Samba"
author      "Guillaume Boudreau"

start on started mysql and started smbd
stop on shutdown

exec ionice -c 2 -n 7 /usr/bin/greyhole --daemon
