# DEBIAN/postinst
#!/bin/bash
set -e

mkdir -p /var/spool/greyhole
chmod 777 /var/spool/greyhole

if [ -d /usr/lib/x86_64-linux-gnu ]; then 
    SOURCE_LIBDIR="/usr/lib64" # Makefile will always install Greyhole .so files in /usr/lib64, for x86_64 CPUs. @see Makefile
    TARGET_LIBDIR="/usr/lib/x86_64-linux-gnu"
elif [ "`uname -m`" = "x86_64" ]; then
    SOURCE_LIBDIR="/usr/lib64"
    TARGET_LIBDIR="/usr/lib64"
    
    # For Ubuntu 11.10, where even x86_64 install use /usr/lib
    if [ -d /usr/lib/samba/vfs/ ]; then
        TARGET_LIBDIR="/usr/lib"
    fi
else
    SOURCE_LIBDIR="/usr/lib"
    TARGET_LIBDIR="/usr/lib"
fi

TARGET_SYMLINK="${TARGET_LIBDIR}/samba/vfs/greyhole.so"

SMB_VERSION="`smbd --version | awk '{print $2}' | awk -F'-' '{print $1}' | awk -F'.' '{print $1,$2}'`"
if [ "${SMB_VERSION}" = "3 4" ]; then
    LIB_FILE="greyhole-samba34.so"
elif [ "${SMB_VERSION}" = "3 5" ]; then
    LIB_FILE="greyhole-samba35.so"
elif [ "${SMB_VERSION}" = "3 6" ]; then
    LIB_FILE="greyhole-samba36.so"
elif [ "${SMB_VERSION}" = "4 0" ]; then
    LIB_FILE="greyhole-samba40.so"
elif [ "${SMB_VERSION}" = "4 1" -a "`uname -m`" = "x86_64" ]; then
    LIB_FILE="greyhole-samba41.so"
else
	echo "Warning: Greyhole doesn't include a VFS module for your version of Samba (${SMB_VERSION})."
	echo "We will try to use the VFS for Samba 3.6, but that might not work."
    LIB_FILE="greyhole-samba36.so"
fi

SOURCE_LIB="${SOURCE_LIBDIR}/greyhole/${LIB_FILE}"

rm -f ${TARGET_SYMLINK}
ln -s ${SOURCE_LIB} ${TARGET_SYMLINK}

if [ -f /proc/fs/cifs/OplockEnabled ]; then
	# cifs client workaround
	# Ref: http://blog.dhampir.no/content/cifs-vfs-no-response-for-cmd-n-mid
	modprobe cifs
	echo 0 > /proc/fs/cifs/OplockEnabled
fi
if [ -f /sys/module/cifs/parameters/enable_oplocks ]; then
	# cifs client workaround
	# Ref: http://blog.dhampir.no/content/cifs-vfs-no-response-for-cmd-n-mid
	modprobe cifs enable_oplocks=0
	echo 0 > /sys/module/cifs/parameters/enable_oplocks
fi

running=0
if [ -f /sbin/start ]; then
	# Using Upstart instead of SYSV init.d
	rm -f /etc/init.d/greyhole
	if [ -f /etc/init/greyhole.conf ]; then
		if [ "`status greyhole | grep running | wc -l`" = "1" ]; then
			restart greyhole
			running=1
		fi
	fi
else
	# (SYSV) Service install & start
	rm -f /etc/init/greyhole.conf
	if [ -f /etc/init.d/mysql ]; then
	    sed -i 's/mysqld/mysql/' /etc/init.d/greyhole
    fi
    if [ -f /etc/init.d/samba ]; then
        sed -i 's/smb/samba/' /etc/init.d/greyhole
    fi
	if [ -f /sbin/chkconfig ]; then
		/sbin/chkconfig --add greyhole
		/sbin/chkconfig greyhole on
	else
		/usr/sbin/update-rc.d greyhole defaults
	fi
	if [ -f /etc/init.d/greyhole ]; then
		if [ "`service greyhole status | grep 'is running' | wc -l`" = "1" ]; then
			service greyhole restart
			running=1
		fi
	fi
fi

if [ $running -eq 0 ]; then
	echo "==========================================================================="
	echo "See /usr/share/greyhole/USAGE to learn how to configure and start Greyhole."
	echo "==========================================================================="
fi
